# social/kamiyo-scheduler.service
[Unit]
Description=Kamiyo AI Agent Quote Tweet Random Scheduler
After=network.target

[Service]
Type=simple
User=kamiyo
WorkingDirectory=/opt/kamiyo
ExecStart=/usr/bin/python3 /opt/kamiyo/social/random_scheduler.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/kamiyo/scheduler.log
StandardError=append:/var/log/kamiyo/scheduler.error.log

# Environment variables (copy from your .env file)
EnvironmentFile=/opt/kamiyo/.env

[Install]
WantedBy=multi-user.target
